# yaasr/terminal/tpl/supervisor.ini
# This file is automatically generated, do not edit manually. 

[program:yaasr-{{ stream_name }}]

command=yaasr compress-and-google-store --stream {{ stream_name }} --bucket_name {{ bucket_name }} --google-credentials {{ google_credentials_path }} --total_seconds {{ total_seconds }} --chunk_bytes_size {{ chunk_bytes_size }} --chunk_time_size {{ chunk_time_size }}

directory=/home/{{ user_name }}/yaasr
user={{ user_name }}
autostart=true
autorestart=true
stdout_logfile=/var/log/yaasr-{{ stream_name }}.log
stderr_logfile=/var/log/yaasr-{{ stream_name }}.log